# pi-setup/finance-display.service
[Unit]
Description=Finance Display Server
After=network.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/finance-display
ExecStart=/home/pi/finance-display/server
Restart=on-failure
RestartSec=5
Environment=PORT=3000

[Install]
WantedBy=multi-user.target
